# debian/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        # Warn if COSMIC desktop is not detected
        if [ -n "$XDG_CURRENT_DESKTOP" ] && ! echo "$XDG_CURRENT_DESKTOP" | grep -qi "cosmic"; then
            echo "WARNING: Flux requires the COSMIC desktop environment."
            echo "The current desktop is: $XDG_CURRENT_DESKTOP"
            echo "The package has been installed but will not function without COSMIC."
        fi
        # Update icon cache
        if command -v gtk-update-icon-cache >/dev/null 2>&1; then
            gtk-update-icon-cache -f -t /usr/share/icons/hicolor || true
        fi
        if command -v systemctl >/dev/null 2>&1; then
            # Enable for all users; takes effect at next login and needs no
            # user bus, so it also works during image builds.
            systemctl --global enable cosmic-ext-flux-daemon || true

            # Best-effort immediate handoff for the installing user.
            # systemctl --user needs the user's runtime dir and session bus,
            # which sudo/su do not provide — point at them explicitly.
            uid="${SUDO_UID:-$PKEXEC_UID}"
            if [ -n "$uid" ] && [ -d "/run/user/$uid" ] && command -v runuser >/dev/null 2>&1; then
                user_systemctl() {
                    runuser -u "$(id -nu "$uid")" -- env \
                        XDG_RUNTIME_DIR="/run/user/$uid" \
                        DBUS_SESSION_BUS_ADDRESS="unix:path=/run/user/$uid/bus" \
                        systemctl --user "$@" 2>/dev/null || true
                }
                user_systemctl daemon-reload
                # Stop the pre-rename daemon if it is still running
                user_systemctl stop cosmic-flux-daemon
                user_systemctl disable cosmic-flux-daemon
                user_systemctl start cosmic-ext-flux-daemon
            fi
        fi
        ;;
esac

#DEBHELPER#
