# debian/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        # Warn if COSMIC desktop is not detected
        if [ -n "$XDG_CURRENT_DESKTOP" ] && ! echo "$XDG_CURRENT_DESKTOP" | grep -qi "cosmic"; then
            echo "WARNING: Flux requires the COSMIC desktop environment."
            echo "The current desktop is: $XDG_CURRENT_DESKTOP"
            echo "The package has been installed but will not function without COSMIC."
        fi
        # Update icon cache
        if command -v gtk-update-icon-cache >/dev/null 2>&1; then
            gtk-update-icon-cache -f -t /usr/share/icons/hicolor || true
        fi
        # Enable and start the daemon for the installing user
        if command -v systemctl >/dev/null 2>&1; then
            systemctl --global daemon-reload || true
            # If running as a logged-in user (not during image build), enable the service
            if [ -n "$SUDO_USER" ]; then
                # Clean up the pre-rename service if it lingers, then enable the new one
                su - "$SUDO_USER" -c "systemctl --user disable --now cosmic-flux-daemon 2>/dev/null" || true
                su - "$SUDO_USER" -c "systemctl --user enable --now cosmic-ext-flux-daemon" || true
            elif [ -n "$USER" ] && [ "$USER" != "root" ]; then
                systemctl --user disable --now cosmic-flux-daemon 2>/dev/null || true
                systemctl --user enable --now cosmic-ext-flux-daemon || true
            fi
        fi
        ;;
esac

#DEBHELPER#
